# service/socialfarm/run
#!/bin/bash 

exec /usr/bin/python /root/SocialFarm/server/server.py 2>&1 >>server.log 

